# files/etc/cron.d/backupcleanup
## THIS FILE IS PUPPETIZED, CHANGES WILL BE OVERWRITTEN EVENTUALLY
## CONTACT A SYSADMIN TO CHANGE SETTINGS IN HERE

0 5 * * *	root	/root/scripts/removeBackupTmpFolders.sh
