# files/etc/cron.d/backupcleanup
# PUPPETIZED

0 5 * * *	root	/root/scripts/removeBackupTmpFolders.sh
